# ledscape-colors.service
[Unit]
Description=LEDscape rainbow test
Requires=ledscape.service

[Service]
WorkingDirectory=/home/root/LEDscape/
ExecStart=/home/root/LEDscape/colors.py
KillMode=process

[Install]
WantedBy=multi-user.target
